# amd_gpu_stressTool/daemon/rc.local
#!/bin/bash

set -x
dhclient &
stress_path="/home/test/Desktop/1946_ROCm"

rb=$(sudo shutdown -r +13)
${rb}

source ${stress_path}/ctrl.sh
